# crontab
*/15 * * * * python /usr/src/app/manage.py distribute_reviews
