# crontab
*/15 * * * * root python /usr/src/app/manage.py distribute_reviews >>/var/log/mwrite-peer-review-distribution.log 2>&1
